# OpenWRT-LEDE/package/utils/fan/files/fan.init
#!/bin/sh /etc/rc.common
#  fan init
#

START=99

start() {

	ProcNumber=`ps -ef | grep 'fan' | grep -v 'grep' | wc -l`
	if [ $ProcNumber -eq 0 ]; then
		fan
		echo "==========Run fan APP =========="
	fi
}
