# OpenWRT-LEDE/package/libs/weston/files/weston-init.ini
#!/bin/sh /etc/rc.common

START=99
STOP=00

DAEMON_OPTS='--tty=1'

start() {
        export XDG_RUNTIME_DIR='/tmp'
        start-stop-daemon -m -p /var/run/weston-init.pid -b -S -x /usr/bin/weston -- $DAEMON_OPTS
}


stop() {
        start-stop-daemon -p /var/run/weston-init.pid -K -s SIGINT
}
